# javaserver.service
[Unit]
Description=Java Server

[Service]
Type=simple
ExecStart=/example/server.bash

[Install]
WantedBy=multi-user.target
